# systemd/system/refreshmirrors.service
[Unit]
Description=Script to refresh Pacman's Mirrorlist
After=network.target

[Service]
ExecStart=/usr/local/bin/refreshmirrors
Type=oneshot
